# systemd/headscale/headscale-update.timer
[Unit]
Description=Weekly Headscale auto‑update

[Timer]
OnCalendar=Sun *-*-* 03:30:00     # every Sunday 03:30
Persistent=true                   # catch‑up if the box was off

[Install]
WantedBy=timers.target